# PKGBUILDS/framebufferphone/fbp/fbp.install
post_install() {
	sudo systemctl enable fbp
}

