# PKGBUILDS/framebufferphone/fbp/fbp.service
[Unit]
Description=Starts the fbp UI at boot

[Service]
ExecStart=/usr/sbin/fbp

[Install]
WantedBy=graphical.target
