# PKGBUILDS/framebufferphone/fbp/fbp.install
post_install() {
	# Start fbp and enable fbp on boot
	sudo systemctl enable fbp
	sudo systemctl start fbp

	# Change power key behavior
	sudo sed -i 's/#HandlePowerKey=poweroff/HandlePowerKey=ignore/g' /etc/systemd/logind.conf
}

