# PKGBUILDS/sxmo/sxmo-utils/sxmo-setpermissions.service
[Unit]
Description=sxmo-pinephone
[Service]
ExecStart=/usr/bin/sxmo_setpermissions.sh
[Install]
WantedBy=multi-user.target
