# PKGBUILDS/pine64/actkbd-pinephone-volume/actkbd-pinephone-volume.install
post_install() {
  printf '%s\n' 'To enable volume button control run the following:'
  printf '%s\n' 'sudo systemctl enable --now pinephone-volume-permfix'
  printf '%s\n' 'systemctl enable --now --user actkbd'
}
