# PKGBUILDS/pine64/acpi-kbd-autolock-pinetab2/acpi-kbd-autolock-pinetab2.install
post_install() {
    systemctl enable --now acpid
}
