# PKGBUILDS/plasma/hfd-service/hfd-service.service
[Unit]
Description=HFD Service
Wants=dbus.service

[Service]
ExecStart=/usr/bin/hfd-service

[Install]
WantedBy=multi-user.target
