# PKGBUILDS/danctnix/mmsd-tng/mmsd-tng.install
post_install() {
  systemctl --global enable mmsd-tng.service
}

pre_remove() {
  systemctl --global disable mmsd-tng.service
}
