# PKGBUILDS/plasma/plasma-mobile-settings/90-plasma-flashlight.rules
SUBSYSTEM=="leds", ACTION=="add", KERNEL=="white:flash", \
  RUN+="/bin/chmod 666 /sys/class/leds/%k/brightness"
