# PKGBUILDS/plasma/plasma-mobile-settings/plasma-mobile-settings.install
post_upgrade() {
    if [[ -e /etc/systemd/system/multi-user.target.wants/ofono.service ||
         -e /usr/lib/systemd/system/multi-user.target.wants/ofono.service ]]
    then
        systemctl disable ofono.service
    fi

    systemctl enable ModemManager.service
}

post_install() {
    post_upgrade
}
