# PKGBUILDS/phosh/flashlight/90-flashlight.rules
SUBSYSTEM=="leds", ACTION=="add", KERNEL=="white:flash", \
  RUN+="/bin/chmod 777 /sys/class/leds/%k/brightness"
