# PKGBUILDS/pine64/ppkb-tools/ppkb-usb-inputd.service
[Unit]
Description=PinePhone Keyboard USB Input Daemon

[Service]
Type=simple
ExecStart=/usr/bin/ppkb-usb-inputd
Restart=on-failure

[Install]
WantedBy=sysinit.target
