# PKGBUILDS/danctnix/danctnix-tweaks-app/danctnix-tweaks-app.install
post_install() {
  systemctl enable danctnix-tweakd.service
}

post_remove() {
  systemctl disable danctnix-tweakd.service
}
