# PKGBUILDS/pine64/device-pine64-pinetab/pinetab-camera-setup.service
[Unit]
Description=Configure PinePhone Camera

[Service]
Type=oneshot
ExecStart=/usr/bin/pinephone-camera-setup.sh

[Install]
WantedBy=multi-user.target
