# PKGBUILDS/danctnix/vvmd/vvmd.install
post_install() {
  systemctl --global enable vvmd.service
}

pre_remove() {
  systemctl --global disable vvmd.service
}
