# PKGBUILDS/pine64/device-pine64-pinephonepro/device-pine64-pinephonepro.install
post_install() {
	systemctl enable eg25-manager
}

post_upgrade() {
	systemctl enable eg25-manager
}
