# PKGBUILDS/pine64/device-pine64-pinephonepro/device-pine64-pinephonepro.install
gen_uboot_script() {
    echo "Generating U-Boot script"
    mkimage -A arm -O linux -T script -C none -n "U-Boot boot script" -d /boot/boot.txt /boot/boot.scr
}

post_install() {
    gen_uboot_script
    systemctl enable eg25-manager
}

post_upgrade() {
    gen_uboot_script
    systemctl enable eg25-manager
}
